# supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=0

[program:fastapi]
command=python3 -m uvicorn shortlist.api.app:create_app --factory --host 127.0.0.1 --port 8001 --limit-max-requests 10000
directory=/app
priority=10
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:scheduler]
command=python3 -m shortlist.scheduler
directory=/app
priority=15
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nextjs]
command=bash -c "until curl -sf http://127.0.0.1:8001/api/health; do sleep 1; done; exec node --max-old-space-size=192 /app/web/server.js"
directory=/app/web
environment=PORT="3000",HOSTNAME="0.0.0.0"
priority=20
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
